# deploy/systemd/georoute@.timer
[Unit]
Description=Periodic refresh of %I country prefix feed (georoute)
Documentation=https://github.com/dantte-lp/georoute

[Timer]
# First run shortly after boot, then every 12h with a random spread
# to avoid hammering the upstream feed source at the same instant.
OnBootSec=5min
OnUnitActiveSec=12h
RandomizedDelaySec=30min
AccuracySec=1min
Persistent=true
Unit=georoute@%i.service

[Install]
WantedBy=timers.target
